# exaim.ini
[program:django]
command=python /path/to/your/manage.py runserver 0.0.0.0:8000
directory=/path/to/your/project
autostart=true
autorestart=true
stderr_logfile=./django.err.log
stdout_logfile=./django.out.log
